# deploy/tvbox-cec.service
# tvbox CEC remote -> uinput bridge - systemd USER unit (no root: /dev/uinput
# and /dev/cec* access comes from the udev rule + input/video group membership
# set up by provision.sh). Installed to ~/.config/systemd/user/ by deploy.sh;
# starts at boot via `loginctl enable-linger`.
[Unit]
Description=tvbox CEC remote -> uinput bridge

[Service]
ExecStart=/usr/bin/python3 %h/.tvbox/cec_uinput_bridge.py
Restart=always
RestartSec=3
# Free hardening: the bridge never needs to gain privileges (device access is
# via group membership; poweroff/standby go through logind D-Bus, not setuid).
# Heavier sandboxing (ProtectSystem etc.) would break the shared /tmp FIFOs, so
# keep it to this one knob.
NoNewPrivileges=yes

[Install]
WantedBy=default.target
